# catkin_ws/devel/etc/catkin/profile.d/06-ctr-nuke.sh
version https://git-lfs.github.com/spec/v1
oid sha256:ec8e1690730f39c866f092162efaa8560e041a9944a3adb1894a5b972e5cb87c
size 75
